# styleantik/etc/crontab
# logrotate
00 00 * * *	/bin/sh -c "/usr/local/sbin/logrotate -s /var/run/logrotate.state /etc/logrotate.conf" >/dev/null 2>&1 
